# projects/liars_poker_bot/crates/xtask/euchre-server.service
[Unit]
Description=Euchre Server

[Service]
ExecStart=/root/deploy/euchre_server
WorkingDirectory=/root/deploy/

[Install]
WantedBy=multi-user.target